# services.conf
[supervisord]
nodaemon=true

[program:flask]
command=/opt/venv/bin/python3 /testrunnerapp/app.py
user=user
autostart=true
autorestart=true
